# supervisord.conf
[program:node]
command = /home/dotcloud/current/run
directory = /home/dotcloud/current